# openwrt-builder/targets/example/files/etc/uci-defaults/99-disable-lan-led
#!/bin/sh

[ -f "/etc/.upgrade" ] && exit 0

uci -q batch <<EOF || true
    del system.led_lan
    commit system
EOF

exit 0
