# openwrt-builder/shared/files/etc/uci-defaults/99-bash
#!/bin/sh

if [ -x "/bin/bash" ] && grep -q "/bin/ash" /etc/passwd; then
	sed -i "s,/bin/ash,/bin/bash," /etc/passwd
fi

exit 0
