# openwrt-builder/shared/files/etc/uci-defaults/99-optimism-luci-webui
#!/bin/sh

if [ -f "/etc/uwsgi/vassals/luci-webui.ini" ]; then
	sed -i "/^threads =/c\threads = 2" /etc/uwsgi/vassals/luci-webui.ini
fi

exit 0
